# mongod.conf
[program:mongod]
command=/usr/bin/mongod --fork --logpath /var/log/mongodb.log
autostart=true
autorestart=true
