# oacis_jupyter/jupyter.conf
[program:jupyter]
command=/home/oacis/miniconda3/bin/jupyter notebook
user=oacis
autostart=true
autorestart=true

