# oacis_jupyter/jupyter.conf
[program:jupyter]
command=bash -l -c 'pipenv run jupyter notebook'
directory=/home/oacis
user=oacis
autostart=true
autorestart=true

